# systemd/blog.target
[Unit]
Description=Blog local stack target
Wants=blog-backend.service blog-cloudflared.service
After=network-online.target

[Install]
WantedBy=default.target
