# systemd/blog-cloudflared.service
[Unit]
Description=Blog Cloudflare Tunnel
After=network-online.target blog-backend.service
Wants=network-online.target blog-backend.service
PartOf=blog.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared --config /home/nhatbang/.config/blog/cloudflared.yml tunnel run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
